# arp-minutely-scanner.timer
[Timer]
OnCalendar=*:0/1
RemainAfterElapse=no

[Install]
WantedBy=timers.target
