# arp-minutely-scanner.service
[Unit]
Description=Print ALL ARP links
Documentation=man:arp

[Service]
ExecStart=
ExecStart=/usr/bin/arp-minutely-scanner
